# adapters/dlq/systemd/dlq.service
[Unit]
Description=nervous-bus dead-letter queue daemon (persist + retry + quarantine)
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.cache/nervous-bus
ExecStart=%h/.pyenv/versions/3.11.13/bin/python3 %h/projects/nervous-bus/adapters/dlq/dlq.py
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dlq
ReadWritePaths=%h/.cache/nervous-bus

[Install]
WantedBy=default.target
