# adapters/pattern-watchdog/systemd/pattern-watchdog.timer
[Unit]
Description=pattern pipeline health check — every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=5min
Unit=pattern-watchdog.service
AccuracySec=30s

[Install]
WantedBy=timers.target
