# adapters/debug-log-rotate/systemd/debug-log-rotate.timer
[Unit]
Description=Periodic size check for nervous-bus debug.jsonl rotation

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=debug-log-rotate.service
AccuracySec=1min

[Install]
WantedBy=timers.target
